# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:mysql]
command=/usr/bin/mysqld_safe
autostart=true
autorestart=true
user=mysql
stdout_logfile=/var/log/supervisor/mysql.log
stderr_logfile=/var/log/supervisor/mysql_err.log

[program:litecp-agent]
command=/opt/litecp/bin/litecp-agent --socket /opt/litecp/run/agent.sock --log-level debug
autostart=true
autorestart=true
user=root
stdout_logfile=/var/log/supervisor/litecp-agent.log
stderr_logfile=/var/log/supervisor/litecp-agent_err.log

[program:litecp]
command=/opt/litecp/bin/litecp --data-dir /opt/litecp/data --agent-socket /opt/litecp/run/agent.sock --listen :3050 --log-level debug
autostart=true
autorestart=true
user=root
stdout_logfile=/var/log/supervisor/litecp.log
stderr_logfile=/var/log/supervisor/litecp_err.log

[program:caddy]
command=/usr/local/bin/caddy run --config /opt/litecp/config/Caddyfile
autostart=true
autorestart=true
user=root
stdout_logfile=/var/log/supervisor/caddy.log
stderr_logfile=/var/log/supervisor/caddy_err.log
